# docker/supervisord-app.ini
[program:nginx]
command=nginx
autostart=true
autorestart=false
priority=10
startretries=0
stopwaitsecs=1

[program:mgserver]
command=./skynet/skynet ./src/%(ENV_NODETYPE)s/config.lua
autostart=true
autorestart=false
priority=4
startretries=0
stopwaitsecs=1
directory=/app/src
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/mgserver_err.log

[eventlistener:processes]
command=stop-supervisor.sh
events=PROCESS_STATE_STOPPED, PROCESS_STATE_EXITED, PROCESS_STATE_FATAL
